# setup_3g_4g_on_boot.timer
[Unit]
Description=timer unit for enabling network only on boot

[Timer]
OnBootSec=60
Unit=setup_3g_4g.service

[Install]
WantedBy=multi-user.target
